# env/etc/cron.d/magento
* * * * * #user# /bin/bash #documentroot#/cron.sh